# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx512M -Xms512M -jar minecraft_server.jar

