# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx1536M -Xms1536M -jar minecraft_server.jar nogui

