# supervisor/conf.d/minecraft.conf
[program:minecraft]
startsecs=5
directory=/data/
command=java -Xmx256M -Xms256M -jar minecraft_server.jar

