# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx1.5G -XX:MaxPermSize=256M -jar mcpc+.jar nogui
startsecs=60

