# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx768M -Xms768M -jar minecraft_server.jar nogui

