# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx1536G -XX:MaxPermSize=256M -jar cauldron-1.6.4-1.965.4.jar nogui
startsecs=60

