# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx1024M -XX:MaxPermSize=256M -jar minecraft_server.jar nogui

