# supervisor/conf.d/minecraft.conf
[program:minecraft]
directory=/data/
command=java -Xmx1024M -Xms1024M -jar minecraft_server.jar nogui

