# supervisor/conf.d/minecraft.conf
[program:minecraft]
startsecs=5
directory=/data/
command=java -Xmx512M -Xms512M -jar minecraft_server.jar

